# 1.0.Linux.Базовый_курс/Lesson6/my_ping.service
[Unit]
Description=My Ping Service
After=multi-user.target

[Service]
Type = simple
ExecStart=/home/kostya/my_ping

[Install]
WantedBy=multi-user.target